# systemd/proxytor-token-rotate.timer
[Unit]
Description=Rotate ProxyTor admin token every 24 hours

[Timer]
OnBootSec=10min
OnUnitActiveSec=24h
Persistent=true
Unit=proxytor-token-rotate.service

[Install]
WantedBy=timers.target
