# systemd/proxytor-token-rotate.service
[Unit]
Description=Rotate ProxyTor admin token

[Service]
Type=oneshot
ExecStart=/opt/proxytor-api/scripts/rotate-token.sh
